# qbittorrent
#!/sbin/openrc-run
# author: abdul2926

start() {
        ebegin "Starting qBittorrent"
        start-stop-daemon --make-pidfile --pidfile /run/qbittorrent.pid --quiet --background --start qbittorrent-nox
        eend $?
}

stop() {
        ebegin "Stopping qBittorrent"
        start-stop-daemon --stop --verbose --pidfile /run/qbittorrent.pid
        eend $?
}
